# runit/resque-scheduler/run
#!/bin/bash
exec 2>&1

BIN_PATH=/home/myuser/applications/myapp/sv/scheduler/bin

exec $BIN_PATH RAILS_ENV=staging
